# supervisor/megaman-usb-proxy.ini
[program:megaman-usb-proxy]
command=/opt/MegaManFE/feproxy.py
redirect_stdout=true
redirect_stderr=true
stdout_logfile=/var/log/feproxy-stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stderr_logfile=/var/log/feproxy-stderr.log
stderr_ogfile_maxbytes=1MB
stderr_logfile_backups=5
