# supervisor/megaman-frontend.ini
[program:megaman-frontend]
environment=MIDASBASE="/opt/midas",MIDASSYS="/opt/midas",ROOTANASYS="/opt/rootana",JSROOTSYS="/opt/jsroot",MIDAS_EXPT_NAME="exo",MIDAS_EXPTAB="/opt/midas/online/exptab"
command=/opt/MegaManFE/femegaman.py -i 0
stdout_logfile=/var/log/femegaman-stdout.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stderr_logfile=/var/log/femegaman-stderr.log
stderr_ogfile_maxbytes=1MB
stderr_logfile_backups=5
