# supervisor/megaman_usb_proxy.ini
[program:megaman_usb_proxy]
command=/opt/MegaManFE/feproxy.py
